# container-sync.timer
[Unit]
Description=Run Container Sync every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
